# packaging/systemd/gajae-app.service
# Transitional source template kept under the pre-rename unit filename so a
# deployment still running the pre-Gaminus manager can update across the
# rename; that manager renders the template it finds at this exact path inside
# the update candidate. The rendered legacy-named unit is retired by
# scripts/gaminus.sh on its next install or update. Do not install directly.
[Unit]
Description=Gaminus self-hosted web application (legacy unit name)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@APP_ROOT@
Environment=HOST=@HOST@
Environment=SERVER_PORT=@PORT@
ExecStart=@NODE_BIN@ @APP_ROOT@/scripts/gaminus-runtime.mjs start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gaminus

[Install]
WantedBy=default.target
