# packaging/systemd/gaminus.service
# Source template for the gaminus.service user unit; do not install it directly.
# The manager keeps runtime data and staged releases in ~/.gaminus, then replaces every placeholder.
[Unit]
Description=Gaminus self-hosted web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@APP_ROOT@
Environment=HOST=@HOST@
Environment=SERVER_PORT=@PORT@
ExecStart=@NODE_BIN@ @APP_ROOT@/scripts/gaminus-runtime.mjs start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gaminus

[Install]
WantedBy=default.target
